# pacalias.service
[Service]
ExecStart=pacalias -c /etc/pacalias.conf -n
[Install]
WantedBy=multi-user.target